# data/systemd/fragoi-udpnotify.service
[Unit]
Description=UDP Notify

[Service]
Type=exec
ExecStart=udpnotify

[Install]
WantedBy=graphical-session.target
